# services/connect/run
#!/bin/bash

cd /app/src/
exec node $SERVICE_ENTRY
