# services/connect/run
#!/bin/bash

cd /app/src/
node $SERVICE_ENTRY
